# data_record.conf
[program:data_record]
directory=/home/pi/Grafana_DashBoard_for_EMS
autostart=true
autorestart=true
startretries=10
user=pi
logfile_maxbytes=0
logfile_backups=0
command=python app.py
